# rootfs/etc/s6-overlay/s6-rc.d/init-ctbrec-config/up
/etc/s6-overlay/s6-rc.d/init-ctbrec-config/run